# debian/ppmsens-1/debian/ppmsens.service
[Unit]
Description=Gather air quality data
Wants=dev-co2sens.device

[Service]
User=root
ExecStart=/usr/bin/ppmsens
