# debian/tasmonitor-1/debian/tasmonitor/lib/systemd/system/tasmonitor.service
[Unit]
Description=Co2 and temperature monitor
Wants=dev-co2sens.device

[Service]
User=root
ExecStart=/usr/bin/co2sens
